# applications/roles/node-exporter/templates/node_exporter.service.j2
[Unit]
Description=Node Exporter {{  version }}
After=network-online.target

[Service]
User={{ user }}
Group={{ group }}
Type=simple
ExecStart={{ path }}

[Install]
WantedBy=multi-user.target
